# telemetry.service
[Unit]
Description=Telemetry
Requires=mariadb.service

[Service]
ExecStart=/opt/telemetry/run.sh

[Install]
WantedBy=multi-user.target